# launchd/backup-pre-namespace-2026-05-23/ai.openclaw.qdrant-native.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.openclaw.qdrant-native</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/chrischo/.local/bin/qdrant</string>
        <string>--config-path</string>
        <string>/Users/chrischo/server/brain/qdrant-config.yaml</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/chrischo/server/brain</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>/Users/chrischo/server/brain/logs/qdrant-native.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/chrischo/server/brain/logs/qdrant-native.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/chrischo</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>8192</integer>
    </dict>
    <key>HardResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>16384</integer>
    </dict>
</dict>
</plist>
