# launchd/ai.brain.neo4j.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.brain.neo4j</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/neo4j</string>
    <string>console</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>JAVA_HOME</key>
    <string>/opt/homebrew/opt/openjdk@17</string>
    <key>HOME</key>
    <string>/Users/chrischo</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/chrischo/server/brain/logs/neo4j-native.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/chrischo/server/brain/logs/neo4j-native.err.log</string>
  <key>WorkingDirectory</key>
  <string>/opt/homebrew/var/neo4j</string>
</dict>
</plist>
