# launchd/ai.brain.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>BRAIN_ATOMS_ENABLED</key>
		<string>true</string>
		<key>BRAIN_ATOMS_READ</key>
		<string>true</string>
		<key>BRAIN_AUTO_HEAL_ENABLED</key>
		<string>true</string>
		<key>BRAIN_AUTO_HEAL_SOURCES</key>
		<string>slo_monitor</string>
		<key>BRAIN_CE_MPS_EMPTY_CACHE</key>
		<string>true</string>
		<key>BRAIN_COMMUNITY_SUMMARIES</key>
		<string>1</string>
		<key>BRAIN_COMMUNITY_SUMMARIES_MAX_COMMUNITIES</key>
		<string>3</string>
		<key>BRAIN_COMMUNITY_SUMMARIES_MAX_ENTITIES</key>
		<string>40</string>
		<key>BRAIN_COMMUNITY_SUMMARIES_TIMEOUT_S</key>
		<string>20</string>
		<key>BRAIN_CROSS_ENCODER_ADAPTIVE</key>
		<string>true</string>
		<key>BRAIN_CROSS_ENCODER_ENABLED</key>
		<string>true</string>
		<key>BRAIN_CROSS_ENCODER_SIMPLE_TOP_K</key>
		<string>10</string>
		<key>BRAIN_CROSS_ENCODER_TOP_K</key>
		<string>14</string>
		<key>BRAIN_EMBED_MODEL</key>
		<string>blaifa/multilingual-e5-large-instruct</string>
		<key>BRAIN_EPISODIC_BINDING_ENABLED</key>
		<string>true</string>
		<key>BRAIN_LOOP_AGENT_DISPATCH_DISABLED</key>
		<string>true</string>
		<key>BRAIN_LTR_ENABLED</key>
		<string>true</string>
		<key>BRAIN_MMR_DIVERSITY_ENABLED</key>
		<string>true</string>
		<key>BRAIN_ONTOLOGY_CONDITIONAL_EXPANSION_ENABLED</key>
		<string>true</string>
		<key>BRAIN_ONTOLOGY_EXPANSION_ENABLED</key>
		<string>true</string>
		<key>BRAIN_ONTOLOGY_EXPANSION_MAX_TERMS</key>
		<string>5</string>
		<key>BRAIN_ONTOLOGY_EXPANSION_MODE</key>
		<string>sidecar</string>
		<key>BRAIN_ONTOLOGY_EXPANSION_RELATIONS</key>
		<string>has_agent,owned_by,owns,proxies,depends_on,manages,has_event,prefers</string>
		<key>BRAIN_ONTOLOGY_EXPANSION_SOURCE</key>
		<string>neo4j</string>
		<key>BRAIN_ONTOLOGY_SIDECAR_LIMIT</key>
		<string>2</string>
		<key>BRAIN_PARENT_CHILD_EXPAND</key>
		<string>1</string>
		<key>BRAIN_RERANKER_MODE</key>
		<string>worker</string>
		<key>BRAIN_RERANKER_TIMEOUT_MS</key>
		<string>1000</string>
		<key>BRAIN_RERANKER_URL</key>
		<string>http://127.0.0.1:8792</string>
		<key>BRAIN_SALIENCE_RANKING_ENABLED</key>
		<string>true</string>
		<key>BRAIN_SELF_RAG_ENABLED</key>
		<string>true</string>
		<key>BRAIN_SERVER_PORT</key>
		<string>8791</string>
		<key>BRAIN_SPREADING_ACTIVATION_ENABLED</key>
		<string>true</string>
		<key>HOME</key>
		<string>/Users/chrischo</string>
		<key>KUMA_CRED_FILE</key>
		<string>/Users/chrischo/.brain/credentials/uptime-kuma.json</string>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
	</dict>
	<key>KeepAlive</key>
	<dict>
		<key>Crashed</key>
		<true/>
		<key>SuccessfulExit</key>
		<false/>
	</dict>
	<key>Label</key>
	<string>ai.brain.server</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/chrischo/server/brain/.venv/bin/python</string>
		<string>/Users/chrischo/server/brain/server.py</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>SoftResourceLimits</key>
	<dict>
		<key>NumberOfFiles</key>
		<integer>8192</integer>
	</dict>
	<key>HardResourceLimits</key>
	<dict>
		<key>NumberOfFiles</key>
		<integer>16384</integer>
	</dict>
	<key>StandardErrorPath</key>
	<string>/Users/chrischo/server/brain/logs/server.err.log</string>
	<key>StandardOutPath</key>
	<string>/Users/chrischo/server/brain/logs/server.log</string>
	<key>ThrottleInterval</key>
	<integer>10</integer>
</dict>
</plist>
