# launchd/ai.brain.qdrant-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.brain.qdrant-backup</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/chrischo/server/brain/.venv/bin/python</string>
        <string>/Users/chrischo/server/brain/cli/backup_qdrant.py</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/chrischo/server/brain/logs/qdrant-backup.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/chrischo/server/brain/logs/qdrant-backup.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/chrischo</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>QDRANT_URL</key>
        <string>http://127.0.0.1:6333</string>
        <key>VECTOR_BACKEND</key>
        <string>qdrant</string>
    </dict>
</dict>
</plist>
