# launchd/ai.openclaw.chromadb-native.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.openclaw.chromadb-native</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/chrischo/server/rag/chroma-venv/bin/chroma</string>
    <string>run</string>
    <string>--host</string>
    <string>127.0.0.1</string>
    <string>--port</string>
    <string>8000</string>
    <string>--path</string>
    <string>/Users/chrischo/server/rag/chroma-data</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>IS_PERSISTENT</key>
    <string>TRUE</string>
    <key>ANONYMIZED_TELEMETRY</key>
    <string>FALSE</string>
    <key>HOME</key>
    <string>/Users/chrischo</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/chrischo/server/brain/logs/chromadb-native.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/chrischo/server/brain/logs/chromadb-native.err.log</string>
</dict>
</plist>
